# lib/systemd/system/uniloader.path
[Unit]
Description=Uniloader queuemetric service change config file

[Path]
PathChanged=/etc/uniloader/uniloader.json
